# aws/ec2/percycraft.service
[Unit]
Description=Percycraft cluster
Requires=docker.service
After=docker.service

[Service]
Restart=always
ExecStart=docker-compose -f /percycraft/docker-compose.mc.yml --detach --env-file /percycraft/.env up
ExecStop=docker stop $(docker ps -a -q)
TimeoutSec=1200

[Install]
WantedBy=multi-user.target